# ubuntu-ebf-rk3588/overlay/etc/systemd/system/plymouth-quit-wait.service.d/check-hdmi.conf
[Service]
ExecStartPre=/usr/local/bin/check-hdmi-plymouth.sh
TimeoutStartSec=5
